# docker/supervisor_config/hcc-gdrive.conf
; supervisor config for the hcc-photo-upload Google Drive uploader

[program:hcc_gdrive_upload]
command=hcc-upload-gdrive /srv/hcc-photo-upload/production.ini
