# docker/supervisor_config/hcc.conf
; supervisor config for the hcc-photo-upload site (uWSGI daemon)

[program:hccphotoupload]
command=/usr/bin/uwsgi --ini-paste /srv/hcc-photo-upload/production.ini
