# docker/supervisor_config/nginx.conf
; supervisor config for nginx

[program:nginx]
command=/usr/sbin/nginx
